# conf/tags.yml
# ------------------------------------------------------------------
# Recursive-IR configuration file
# Copyright (c) 2026 Mark Jayson Alvarez
# Licensed under the Recursive-IR License
# ------------------------------------------------------------------

taxonomies:
  - id: "MITRE"
    name: "MITRE ATT&CK"
    description: "Enterprise ATT&CK tactics, techniques, and subtechniques for tagging events."
    tactics:

      - id: "TA0043"
        name: "Reconnaissance"
        tag: "MITRE-TA0043-Reconnaissance"
        description: "Adversary gathering information for planning and targeting."
        aliases:
          - "recon"
          - "information gathering"
          - "footprinting"
          - "scanning"
        techniques:
          - id: "T1595"
            name: "Active Scanning"
            tag: "MITRE-T1595-Active-Scanning"
            description: "Scanning for live hosts, services, and vulnerabilities."
            subtechniques:
              - id: "T1595.001"
                name: "Scanning IP Blocks"
              - id: "T1595.002"
                name: "Vulnerability Scanning"

          - id: "T1590"
            name: "Gather Victim Network Information"
            tag: "MITRE-T1590-Gather-Victim-Network-Information"
            description: "Collecting info about victim network topology and infrastructure."
            subtechniques:
              - id: "T1590.001"
                name: "Domain Properties"
              - id: "T1590.002"
                name: "DNS"
              - id: "T1590.006"
                name: "Network Security Appliances"

          - id: "T1589"
            name: "Gather Victim Identity Information"
            tag: "MITRE-T1589-Gather-Victim-Identity-Information"
            description: "Collecting info on accounts and identities."
            subtechniques:
              - id: "T1589.001"
                name: "Credentials"
              - id: "T1589.002"
                name: "Email Addresses"

          - id: "T1598"
            name: "Phishing for Information"
            tag: "MITRE-T1598-Phishing-for-Information"
            description: "Social engineering to elicit data from targets."
            subtechniques:
              - id: "T1598.001"
                name: "Spearphishing Service"
              - id: "T1598.002"
                name: "Spearphishing Link"

      - id: "TA0042"
        name: "Resource Development"
        tag: "MITRE-TA0042-Resource-Development"
        description: "Establishing resources to support operations (infra, tooling, accounts)."
        aliases:
          - "infrastructure setup"
          - "tool acquisition"
        techniques:
          - id: "T1583"
            name: "Acquire Infrastructure"
            tag: "MITRE-T1583-Acquire-Infrastructure"
            description: "Obtaining servers, domains, and network resources."
            subtechniques:
              - id: "T1583.001"
                name: "Domains"
              - id: "T1583.002"
                name: "DNS Servers"
              - id: "T1583.003"
                name: "Virtual Private Servers"

          - id: "T1584"
            name: "Compromise Infrastructure"
            tag: "MITRE-T1584-Compromise-Infrastructure"
            description: "Taking over existing infrastructure for malicious use."
            subtechniques:
              - id: "T1584.001"
                name: "Domains"
              - id: "T1584.002"
                name: "DNS Servers"

          - id: "T1608"
            name: "Stage Capabilities"
            tag: "MITRE-T1608-Stage-Capabilities"
            description: "Staging malware, tools, and payloads for deployment."
            subtechniques:
              - id: "T1608.001"
                name: "Upload Malware"
              - id: "T1608.002"
                name: "Install Digital Certificate"

      - id: "TA0001"
        name: "Initial Access"
        tag: "MITRE-TA0001-Initial-Access"
        description: "Gaining a foothold in the environment."
        aliases:
          - "initial compromise"
          - "entry point"
        techniques:
          - id: "T1566"
            name: "Phishing"
            tag: "MITRE-T1566-Phishing"
            description: "Social engineering delivered via email or messaging."
            subtechniques:
              - id: "T1566.001"
                name: "Spearphishing Attachment"
              - id: "T1566.002"
                name: "Spearphishing Link"
              - id: "T1566.003"
                name: "Spearphishing via Service"

          - id: "T1190"
            name: "Exploit Public-Facing Application"
            tag: "MITRE-T1190-Exploit-Public-Facing-Application"
            description: "Exploitation of Internet-exposed applications."

          - id: "T1133"
            name: "External Remote Services"
            tag: "MITRE-T1133-External-Remote-Services"
            description: "Use of VPN, RDP, or other remote access for entry."

          - id: "T1078"
            name: "Valid Accounts"
            tag: "MITRE-T1078-Valid-Accounts"
            description: "Using stolen or legitimate credentials for access."
            subtechniques:
              - id: "T1078.001"
                name: "Local Accounts"
              - id: "T1078.002"
                name: "Domain Accounts"
              - id: "T1078.004"
                name: "Cloud Accounts"

      - id: "TA0002"
        name: "Execution"
        tag: "MITRE-TA0002-Execution"
        description: "Running malicious or attacker-controlled code."
        aliases:
          - "code execution"
          - "malware execution"
        techniques:
          - id: "T1059"
            name: "Command and Scripting Interpreter"
            tag: "MITRE-T1059-Command-and-Scripting-Interpreter"
            description: "Execution via shells and scripting engines."
            subtechniques:
              - id: "T1059.001"
                name: "PowerShell"
              - id: "T1059.003"
                name: "Windows Command Shell"
              - id: "T1059.005"
                name: "Visual Basic"
              - id: "T1059.006"
                name: "Python"

          - id: "T1204"
            name: "User Execution"
            tag: "MITRE-T1204-User-Execution"
            description: "Users running content that leads to execution."
            subtechniques:
              - id: "T1204.002"
                name: "Malicious File"
              - id: "T1204.001"
                name: "Malicious Link"

          - id: "T1105"
            name: "Ingress Tool Transfer"
            tag: "MITRE-T1105-Ingress-Tool-Transfer"
            description: "Transferring tools or payloads onto a compromised system."

          - id: "T1047"
            name: "Windows Management Instrumentation"
            tag: "MITRE-T1047-Windows-Management-Instrumentation"
            description: "Use of WMI for remote or local execution."

      - id: "TA0003"
        name: "Persistence"
        tag: "MITRE-TA0003-Persistence"
        description: "Maintaining access and foothold across restarts and logoffs."
        aliases:
          - "backdoor"
          - "autorun persistence"
        techniques:
          - id: "T1053"
            name: "Scheduled Task/Job"
            tag: "MITRE-T1053-Scheduled-Task-Job"
            description: "Using scheduled tasks or jobs to run payloads."
            subtechniques:
              - id: "T1053.005"
                name: "Scheduled Task (Windows)"
              - id: "T1053.003"
                name: "Cron"

          - id: "T1060"
            name: "Registry Run Keys / Startup Folder"
            tag: "MITRE-T1060-Registry-Run-Keys-Startup-Folder"
            description: "Persisting via startup locations on Windows."

          - id: "T1547"
            name: "Boot or Logon Autostart Execution"
            tag: "MITRE-T1547-Boot-or-Logon-Autostart-Execution"
            description: "Mechanisms that run automatically at boot or logon."
            subtechniques:
              - id: "T1547.001"
                name: "Registry Run Keys"
              - id: "T1547.009"
                name: "Shortcut Modification"

          - id: "T1136"
            name: "Create Account"
            tag: "MITRE-T1136-Create-Account"
            description: "Creating local or domain accounts for long-term access."
            subtechniques:
              - id: "T1136.001"
                name: "Local Account"
              - id: "T1136.002"
                name: "Domain Account"

      - id: "TA0004"
        name: "Privilege Escalation"
        tag: "MITRE-TA0004-Privilege-Escalation"
        description: "Gaining higher-level permissions."
        aliases:
          - "privilege elevation"
        techniques:
          - id: "T1068"
            name: "Exploitation for Privilege Escalation"
            tag: "MITRE-T1068-Exploitation-for-Privilege-Escalation"
            description: "Abusing software flaws to elevate privileges."

          - id: "T1134"
            name: "Access Token Manipulation"
            tag: "MITRE-T1134-Access-Token-Manipulation"
            description: "Manipulating tokens for privilege/identity changes."
            subtechniques:
              - id: "T1134.001"
                name: "Token Impersonation/Theft"

          - id: "T1078"
            name: "Valid Accounts"
            tag: "MITRE-T1078-Valid-Accounts"
            description: "Reusing higher-privileged accounts to escalate."

          - id: "T1548"
            name: "Abuse Elevation Control Mechanism"
            tag: "MITRE-T1548-Abuse-Elevation-Control-Mechanism"
            description: "Abusing sudo, UAC, or similar."
            subtechniques:
              - id: "T1548.002"
                name: "Bypass User Account Control"

      - id: "TA0005"
        name: "Defense Evasion"
        tag: "MITRE-TA0005-Defense-Evasion"
        description: "Avoiding detection and bypassing defenses."
        aliases:
          - "evasion"
          - "security bypass"
        techniques:
          - id: "T1562"
            name: "Impair Defenses"
            tag: "MITRE-T1562-Impair-Defenses"
            description: "Disabling or weakening security tools."
            subtechniques:
              - id: "T1562.001"
                name: "Disable or Modify Tools"
              - id: "T1562.004"
                name: "Disable System Logging"

          - id: "T1070"
            name: "Indicator Removal on Host"
            tag: "MITRE-T1070-Indicator-Removal-on-Host"
            description: "Removing, altering, or covering logs and artifacts."
            subtechniques:
              - id: "T1070.001"
                name: "Clear Windows Event Logs"
              - id: "T1070.004"
                name: "File Deletion"

          - id: "T1036"
            name: "Masquerading"
            tag: "MITRE-T1036-Masquerading"
            description: "Disguising artifacts as legitimate."

          - id: "T1140"
            name: "Deobfuscate/Decode Files or Information"
            tag: "MITRE-T1140-Deobfuscate-Decode-Files-or-Information"
            description: "Using obfuscation/encoding to hide code or data."

      - id: "TA0006"
        name: "Credential Access"
        tag: "MITRE-TA0006-Credential-Access"
        description: "Stealing account credentials and secrets."
        aliases:
          - "credential dumping"
          - "password theft"
        techniques:
          - id: "T1003"
            name: "OS Credential Dumping"
            tag: "MITRE-T1003-OS-Credential-Dumping"
            description: "Extracting credentials from OS components."
            subtechniques:
              - id: "T1003.001"
                name: "LSASS Memory"
              - id: "T1003.002"
                name: "Security Account Manager (SAM)"
              - id: "T1003.004"
                name: "LSA Secrets"

          - id: "T1555"
            name: "Credentials from Password Stores"
            tag: "MITRE-T1555-Credentials-from-Password-Stores"
            description: "Pulling credentials from browsers or password managers."
            subtechniques:
              - id: "T1555.003"
                name: "Credentials from Web Browsers"

          - id: "T1110"
            name: "Brute Force"
            tag: "MITRE-T1110-Brute-Force"
            description: "Guessing passwords for accounts."
            subtechniques:
              - id: "T1110.001"
                name: "Password Guessing"
              - id: "T1110.003"
                name: "Password Spraying"

          - id: "T1552"
            name: "Unsecured Credentials"
            tag: "MITRE-T1552-Unsecured-Credentials"
            description: "Finding creds in files, scripts, or backups."

      - id: "TA0007"
        name: "Discovery"
        tag: "MITRE-TA0007-Discovery"
        description: "Understanding the environment, accounts, and assets."
        aliases:
          - "environment enumeration"
          - "network discovery"
        techniques:
          - id: "T1087"
            name: "Account Discovery"
            tag: "MITRE-T1087-Account-Discovery"
            description: "Enumerating users and groups."
            subtechniques:
              - id: "T1087.001"
                name: "Local Account"
              - id: "T1087.002"
                name: "Domain Account"

          - id: "T1046"
            name: "Network Service Scanning"
            tag: "MITRE-T1046-Network-Service-Scanning"
            description: "Scanning for open ports and services."

          - id: "T1082"
            name: "System Information Discovery"
            tag: "MITRE-T1082-System-Information-Discovery"
            description: "Querying OS version, hardware, etc."

          - id: "T1482"
            name: "Domain Trust Discovery"
            tag: "MITRE-T1482-Domain-Trust-Discovery"
            description: "Enumerating trust relationships across domains."

      - id: "TA0008"
        name: "Lateral Movement"
        tag: "MITRE-TA0008-Lateral-Movement"
        description: "Moving from one host or account to another."
        aliases:
          - "internal movement"
          - "pivoting"
        techniques:
          - id: "T1021"
            name: "Remote Services"
            tag: "MITRE-T1021-Remote-Services"
            description: "Using remote management protocols for lateral movement."
            subtechniques:
              - id: "T1021.001"
                name: "Remote Desktop Protocol (RDP)"
              - id: "T1021.002"
                name: "SMB/Windows Admin Shares"
              - id: "T1021.006"
                name: "Windows Remote Management (WinRM)"

          - id: "T1077"
            name: "Windows Admin Shares"
            tag: "MITRE-T1077-Windows-Admin-Shares"
            description: "Using admin shares to move and execute."

          - id: "T1550"
            name: "Use of Stolen Credentials"
            tag: "MITRE-T1550-Use-of-Stolen-Credentials"
            description: "Using captured creds for lateral access."

          - id: "T1027"
            name: "Obfuscated Files or Information"
            tag: "MITRE-T1027-Obfuscated-Files-or-Information"
            description: "Hiding payloads used for lateral tooling."

      - id: "TA0009"
        name: "Collection"
        tag: "MITRE-TA0009-Collection"
        description: "Gathering target data for later use or exfiltration."
        aliases:
          - "data collection"
        techniques:
          - id: "T1114"
            name: "Email Collection"
            tag: "MITRE-T1114-Email-Collection"
            description: "Collecting emails from mail clients or servers."
            subtechniques:
              - id: "T1114.001"
                name: "Local Email Collection"
              - id: "T1114.003"
                name: "Cloud Email Collection"

          - id: "T1113"
            name: "Screen Capture"
            tag: "MITRE-T1113-Screen-Capture"
            description: "Capturing desktop screenshots."

          - id: "T1056"
            name: "Input Capture"
            tag: "MITRE-T1056-Input-Capture"
            description: "Capturing keyboard or mouse input."
            subtechniques:
              - id: "T1056.001"
                name: "Keylogging"

          - id: "T1115"
            name: "Clipboard Data"
            tag: "MITRE-T1115-Clipboard-Data"
            description: "Reading or capturing clipboard contents."

      - id: "TA0011"
        name: "Command and Control"
        tag: "MITRE-TA0011-Command-and-Control"
        description: "Maintaining communication with compromised systems."
        aliases:
          - "c2"
          - "beaconing"
        techniques:
          - id: "T1071"
            name: "Application Layer Protocol"
            tag: "MITRE-T1071-Application-Layer-Protocol"
            description: "Using common protocols to blend C2 traffic."
            subtechniques:
              - id: "T1071.001"
                name: "Web Protocols"
              - id: "T1071.004"
                name: "DNS"

          - id: "T1095"
            name: "Non-Application Layer Protocol"
            tag: "MITRE-T1095-Non-Application-Layer-Protocol"
            description: "Using lower-level protocols for C2."

          - id: "T1105"
            name: "Ingress Tool Transfer"
            tag: "MITRE-T1105-Ingress-Tool-Transfer"
            description: "Transferring tools over C2 channels."

          - id: "T1573"
            name: "Encrypted Channel"
            tag: "MITRE-T1573-Encrypted-Channel"
            description: "Encrypting C2 communications."

      - id: "TA0010"
        name: "Exfiltration"
        tag: "MITRE-TA0010-Exfiltration"
        description: "Stealing or leaking data out of the environment."
        aliases:
          - "data theft"
          - "data exfiltration"
        techniques:
          - id: "T1041"
            name: "Exfiltration Over C2 Channel"
            tag: "MITRE-T1041-Exfiltration-Over-C2-Channel"
            description: "Reusing the C2 channel for data theft."

          - id: "T1048"
            name: "Exfiltration Over Alternative Protocol"
            tag: "MITRE-T1048-Exfiltration-Over-Alternative-Protocol"
            description: "Using non-standard protocols for exfil."
            subtechniques:
              - id: "T1048.002"
                name: "Exfiltration Over Asymmetric Encrypted Non-C2 Protocol"

          - id: "T1567"
            name: "Exfiltration Over Web Services"
            tag: "MITRE-T1567-Exfiltration-Over-Web-Services"
            description: "Using cloud or web services to move stolen data."
            subtechniques:
              - id: "T1567.002"
                name: "Exfiltration to Cloud Storage"

      - id: "TA0040"
        name: "Impact"
        tag: "MITRE-TA0040-Impact"
        description: "Disrupting, destroying, or manipulating systems and data."
        aliases:
          - "destruction"
          - "sabotage"
        techniques:
          - id: "T1486"
            name: "Data Encrypted for Impact"
            tag: "MITRE-T1486-Data-Encrypted-for-Impact"
            description: "Encrypting data for ransom or destruction."

          - id: "T1499"
            name: "Endpoint Denial of Service"
            tag: "MITRE-T1499-Endpoint-Denial-of-Service"
            description: "Overloading endpoints to cause outages."

          - id: "T1489"
            name: "Service Stop"
            tag: "MITRE-T1489-Service-Stop"
            description: "Stopping services to impact availability."

          - id: "T1565"
            name: "Data Manipulation"
            tag: "MITRE-T1565-Data-Manipulation"
            description: "Changing data to cause operational or integrity issues."

